# deploy/systemd/codex-session-agent.service
[Unit]
Description=Codex Session Viewer agent daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wulfuser
Group=wulfuser
WorkingDirectory=/home/wulfuser/codex_viewer
EnvironmentFile=-/home/wulfuser/codex_viewer/deploy/systemd/codex-session-viewer.env
ExecStart=/home/wulfuser/codex_viewer/scripts/start-agent-daemon.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target

