# deploy/systemd/codex-session-agent.service
[Unit]
Description=Agent Operations Viewer agent daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hoff
Group=hoff
WorkingDirectory=/home/hoff/codex_viewer
EnvironmentFile=-/home/hoff/codex_viewer/deploy/systemd/codex-session-viewer.env
ExecStart=/home/hoff/codex_viewer/scripts/start-agent-daemon.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
